# red-empresarial-lima-ansible/roles/switching_l2_safe_standby_fase4c2/templates/jhalex-link-watchdog.service.j2
[Unit]
Description=JHALEX Link Failover Watchdog (Fase 4C.2)
After=network-online.target jhalex-swcorelim1-ovs.service jhalex-swcorelim2-ovs.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/jhalex-link-watchdog
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
