# red-empresarial-lima-ansible/roles/switching_l2_safe_standby_fase4c2/templates/jhalex-fase4c2-safe-topology.service.j2
[Unit]
Description=JHALEX Fase 4C.2 Safe Topology Persistence
After=network-online.target openvswitch-switch.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/jhalex-fase4c2-safe-topology.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
